# devops/dap-service-health.service
[Unit]
Description=DAP API Health service

[Service]
Type=simple
User=diel
WorkingDirectory=/home/diel/API-Server
ExecStart=/usr/bin/node /home/diel/API-Server/dist/serviceHealth/main.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
